# development/jdk/profile.d/jdk.csh
#!/bin/csh
setenv JAVA_HOME /usr/lib/java
setenv PATH ${PATH}:${JAVA_HOME}/bin:${JAVA_HOME}/jre/bin
