# development/openjdk8/profile.d/openjdk8.csh
#!/bin/csh
setenv JAVA_HOME /usr/lib/openjdk8
setenv MANPATH ${MANPATH}:${JAVA_HOME}/man
setenv PATH ${PATH}:${JAVA_HOME}/bin:${JAVA_HOME}/jre/bin
