# system/rubygems/profile.d/gem.csh
#!/bin/csh

setenv GEM_PATH /usr/lib/ruby/gems/1.8
