# system/prelink/prelink.cron
#!/bin/sh

renice +19 -p $$ > /dev/null 2>&1

#By default, we'll do quick cron runs of prelink
PRELINK_OPTS="-amRq"
#If you don't want to do this, then use this opts.
#PRELINK_OPTS="-amR"

/usr/sbin/prelink $PRELINK_OPTS
