# development/gcc5/profile.d/gcc5.csh
setenv CC     gcc-5
setenv CPP    cpp-5
setenv CXX    g++-5
setenv AR     gcc-ar-5
setenv NM     gcc-nm-5
setenv RANLIB gcc-ranlib-5

