# system/mlocate/mlocate.cron
#!/bin/sh
ionice -c3 nice -n 19 /usr/libexec/mlocate-run-updatedb
