# system/linuxconsoletools/90-joystick.rules
# Restore any stored calibration for the device
KERNEL=="js*", ACTION=="add", RUN+="/usr/bin/jscal-restore %E{DEVNAME}"
