# network/dbmail/dbmail.cron
#!/bin/sh

# cleanup deleted messages
if [ -x /usr/sbin/dbmail-util ]; then
	/usr/sbin/dbmail-util -ay 2>&1 > /dev/null
fi

