# libraries/dietlibc/profile.d/dietlibc.csh
#!/bin/csh
setenv MANPATH /opt/diet/man:${MANPATH}
