# libraries/dietlibc/profile.d/dietlibc.csh
#!/bin/csh
setenv PATH ${PATH}:/opt/diet/bin
setenv MANPATH ${MANPATH}:/opt/diet/man
