# libraries/dietlibc/profile.d/dietlibc.sh
#!/bin/sh
export MANPATH="/opt/diet/man:$MANPATH"
