# system/tiger/config/tiger.cron
#
# Regular cron jobs for the tiger package
#
# modified to work with Slackware and Dillon's Cron by pyllyukko
#
0 * * * * /bin/test -x /usr/sbin/tigercron && { DEFAULT=/etc/default/tiger ; [ -r "$DEFAULT" ] && . "$DEFAULT" || NICETIGER=10 ; /bin/nice -n$NICETIGER /usr/sbin/tigercron -q ; }
