# development/jbake/profile.d/jbake.sh
#!/bin/sh
export JBAKE_HOME=/usr/share/jbake
export PATH="${PATH}:${JBAKE_HOME}/bin"
