# development/jbake/profile.d/jbake.csh
#!/bin/csh
setenv JBAKE_HOME /usr/share/jbake
setenv PATH ${PATH}:${JBAKE_HOME}/bin
