# system/plan9port/profile.d/plan9port.csh
#!/bin/csh
setenv PLAN9 /opt/plan9
setenv MANPATH ${MANPATH}:$PLAN9/man
setenv PATH ${PATH}:$PLAN9/bin
