# system/plan9port/profile.d/plan9port.csh
#!/bin/csh
setenv PLAN9 /opt/plan9
setenv MANPATH ${MANPATH}:/opt/plan9/man
setenv PATH ${PATH}:/opt/plan9/bin
