# network/squirrelmail/squirrelmail.cron
#!/bin/sh

find /var/spool/squirrelmail/attach/ -type f -mtime +30 -print0 |
	xargs -0 rm -f >/dev/null 2>&1

