# graphics/blender/profile.d/blender.csh
#!/bin/csh
setenv PATH ${PATH}:/opt/blender
