# development/scala/profile.d/scala.sh
#!/bin/sh
export SCALA_HOME=@LIBDIR@/scala
