# development/scala/profile.d/scala.csh
#!/bin/csh
setenv SCALA_HOME @LIBDIR@/scala
setenv MANPATH ${MANPATH}:${SCALA_HOME}/man
setenv PATH ${PATH}:${SCALA_HOME}/bin
