# development/scala/profile.d/scala.csh
#!/bin/csh
setenv SCALA_HOME @LIBDIR@/scala
