# development/liquibase/profile.d/liquibase.sh
#!/bin/sh

export LIQUIBASE_HOME=/opt/liquibase
