# development/liquibase/profile.d/liquibase.csh
#!/bin/csh
setenv LIQUIBASE_HOME /opt/liquibase
