# development/apache-ant/profile.d/apache-ant.csh
#!/bin/csh
setenv ANT_HOME /opt/apache-ant
setenv PATH ${PATH}:${ANT_HOME}/bin
