# development/apache-ant/profile.d/apache-ant.sh
#!/bin/sh
export ANT_HOME=/opt/apache-ant
export PATH=$PATH:$ANT_HOME/bin
