# misc/ibus-skk/profile.d/ibus-skk.csh
[ -d /usr/share/ibus-skk ]
if ($status == 0) then
  setenv IBUS_SKK_PKGDATADIR /usr/share/ibus-skk
endif
