# misc/ibus-skk/profile.d/ibus-skk.sh
if [ -d /usr/share/ibus-skk ]; then
  export IBUS_SKK_PKGDATADIR=/usr/share/ibus-skk
fi

