# system/graphite-web/graphite-web.cron
#!/bin/sh

set -e

if [ -x /usr/bin/graphite-build-search-index ]; then
	/usr/bin/graphite-build-search-index
fi
