# misc/fcitx/profile.d/fcitx.csh
#!/bin/csh

setenv XMODIFIERS "@im=fcitx"
setenv XIM fcitx
setenv XIM_PROGRAM fcitx
setenv GTK_IM_MODULE fcitx
setenv QT_IM_MODULE fcitx
