# libraries/libvdpau-va-gl/profile.d/vdpau-va-gl.sh
#!/bin/sh

VDPAU_DRIVER=va_gl
export VDPAU_DRIVER=va_gl
