# academic/scilab/profile.d/scilab.csh
#!/bin/csh

setenv PATH ${PATH}:/opt/scilab/bin

setenv PKG_CONFIG_PATH ${PKG_CONFIG_PATH}:/opt/scilab/lib/pkgconfig

