# academic/scilab/profile.d/scilab.sh
#!/bin/sh

PATH=${PATH}:/opt/scilab/bin
export PATH

PKG_CONFIG_PATH=${PKG_CONFIG_PATH}:/opt/scilab/lib/pkgconfig
export PKG_CONFIG_PATH

