# academic/scilab/profile.d/scilab.csh
#!/bin/csh
setenv PATH ${PATH}:/opt/scilab/bin
