# academic/scilab/profile.d/scilab.sh
#!/bin/sh
export PATH=${PATH}:/opt/scilab/bin
export PKG_CONFIG_PATH=${PKG_CONFIG_PATH}:/opt/scilab/lib/pkgconfig
