# graphics/brlcad/profile.d/brlcad.csh
#!/bin/csh
# Add path for brlcad:
setenv PATH ${PATH}:/opt/brlcad/bin

