# audio/ladspa_sdk/profile.d/ladspa.sh
#!/bin/sh
LADSPA_PATH=/usr/lib/ladspa ; export LADSPA_PATH
