# audio/ladspa_sdk/profile.d/ladspa.csh
#!/bin/csh
setenv LADSPA_PATH /usr/lib/ladspa
