# system/heirloom-doctools/profile.d/heirloom.csh
#!/bin/csh
setenv HEIRLOOM /opt/heirloom
setenv MANPATH ${MANPATH}:${HEIRLOOM}/man 
setenv PATH ${PATH}:${HEIRLOOM}/bin
