# development/icecream/profile.d/icecream.sh
#!/bin/sh
export PATH=/usr/libexec/icecc/bin:$PATH
