# development/icecream/profile.d/icecream.csh
#!/bin/csh
setenv PATH /usr/libexec/icecc/bin:${PATH}
