# development/icecream/profile.d/icecream.csh
#!/bin/csh
setenv PATH /opt/icecream/bin:${PATH}
