# development/icecream/profile.d/icecream.sh
#!/bin/sh
export PATH=/opt/icecream/bin:$PATH
