# development/groovy/profile.d/groovy.csh
#!/bin/csh
setenv GROOVY_HOME @LIBDIR@/groovy
setenv PATH ${PATH}:${GROOVY_HOME}/bin
