# development/groovy/profile.d/groovy.sh.new
#!/bin/sh
export GROOVY_HOME=@LIBDIR@/groovy
