# development/groovy/profile.d/groovy.csh.new
#!/bin/csh
setenv GROOVY_HOME @LIBDIR@/groovy
