# network/mailman/profile.d/mailman.sh
#!/bin/sh
export PATH="$PATH:/opt/mailman/bin"
