# network/mailman/profile.d/mailman.csh
#!/bin/csh
setenv PATH ${PATH}:/opt/mailman/bin
