# network/exim/contrib/exim.cron
#!/bin/sh

# Exim's spool directory
SPOOL=/var/spool/exim

# Check for spool and the tidydb utility
test -d $SPOOL -a -x /usr/sbin/exim_tidydb || exit

# Tidy up the contents of the hints databases
find $SPOOL/db -name '*.lockfile' -exec basename {} .lockfile \; \
  | xargs -r -n 1 sudo -u exim /usr/sbin/exim_tidydb -t 7d $SPOOL \
  > /dev/null
