# development/adoptopenjdk-openj9/profile.d/adoptopenjdk-openj9.sh
#!/bin/sh
export JAVA_HOME=/usr/local/lib/jdk-@VERSION@
export MANPATH="${MANPATH}:${JAVA_HOME}/man"
export PATH="${PATH}:${JAVA_HOME}/bin"
